# pocketbase.conf
[program:pb]
command = pocketbase serve
directory = /pb
redirect_stderr = true
stdout_logfile = /var/log/pocketbase.log
stdout_logfile_maxbytes = 1MB
stdout_logfile_backups = 0
